# edtwardy-webservices.service
[Unit]
Description=Webservices
After=docker.service network-online.target
Requires=docker.service network-online.target

[Service]
ExecStart=/bin/start-webservices start

[Install]
WantedBy=multi-user.target