# nginx/renew-certificates.service
[Unit]
Description=Renew twardyece.com TLS Certificates

[Service]
Environment=INTERACTIVE=n
ExecStart=/usr/bin/webservices-certbot renew
