# nginx/twardyece.target
[Unit]
Description=All web-based services running
After=networking.target
Requires=networking.target

[Install]
WantedBy=multi-user.target
