# common/containers-up.target
[Unit]
Description=all containers
Conflicts=containers-down.target
After=containers-down.target

[Install]
WantedBy=multi-user.target
