# common/containers-down.target
[Unit]
Description=all containers
Conflicts=containers-up.target
After=containers-up.target
