# common/pre-containers.target
[Unit]
Description=container infrastructure dependencies
After=podman.service
Requires=podman.service
