# common/services-down.target
[Unit]
Description=All services are stopped
Conflicts=services-up.target
After=services-up.target
