# common/services-down.target
[Unit]
Description=all containers
Conflicts=services-up.target
After=services-up.target
